# .system/config/com.work-wiki.slack-daily.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.work-wiki.slack-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__WORK_WIKI__/.system/scripts/slack-ingest.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__HOME__/.claude/logs/slack-ingest-daily.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/.claude/logs/slack-ingest-daily.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>__HOME__/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>WORK_WIKI_DIR</key>
        <string>__WORK_WIKI__</string>
        <key>WORK_WIKI_SYNTH_PROVIDER</key>
        <string>__SYNTH_PROVIDER__</string>
    </dict>
</dict>
</plist>
